# s6/services.d/eventd-daemon/run
#!/usr/bin/with-contenv sh

cd /usr/share/app
exec python daemon.py
