# s6/services.d/eventd-api/run
#!/usr/bin/with-contenv sh

cd /usr/share/app

exec gunicorn eventd.api.wsgi --bind 0.0.0.0:80 --log-file -
